# scripts/larammerce-supervisor.ini
[program:<user>-admin-queue-worker]
process_name=%(program_name)s_%(process_num)02d
command=/opt/cpanel/ea-php80/root/bin/php /home/<user>/<project-path>/artisan queue:listen --queue <project-name>_ADMIN --memory=5000 --timeout=6000
autostart=true
autorestart=true
user=<user>
numprocs=4
redirect_stderr=true
stdout_logfile=/var/log/supervisor/<user>-admin-queue-worker.log

[program:<user>-customer-queue-worker]
process_name=%(program_name)s_%(process_num)02d
command=/opt/cpanel/ea-php80/root/bin/php /home/<user>/<project-path>/artisan queue:listen --queue <user>_CUSTOMER --memory=5000 --timeout=6000
autostart=true
autorestart=true
user=<user>
numprocs=4
redirect_stderr=true
stdout_logfile=/var/log/supervisor/<user>-customer-queue-worker.log

[program:<user>-general-queue-worker]
process_name=%(program_name)s_%(process_num)02d
command=/opt/cpanel/ea-php80/root/bin/php /home/<user>/<project-path>/artisan queue:listen --memory=5000 --timeout=6000
autostart=true
autorestart=true
user=<user>
numprocs=4
redirect_stderr=true
stdout_logfile=/var/log/supervisor/<user>-general-queue-worker.log